# buildroot-external/board/common/rootfs-overlay/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/usr/bin/glass-agent --addr=:8080 --data-dir=/data
::shutdown:/etc/init.d/rcK

