# buildroot-external/rootfs-overlay/usr/lib/systemd/system/glassos-bind.target
[Unit]
Description=GlassOS bind mounts
Documentation=https://github.com/glasslabs/os
Before=network.target

[Install]
WantedBy=multi-user.target

