# buildroot-external/rootfs-overlay/usr/lib/systemd/system/glassos-wifi-provision.service
[Unit]
Description=GlassOS WiFi provisioning
Documentation=https://github.com/glasslabs/os
# Run after system dirs are created and the connections bind-mount is up,
# but strictly before NetworkManager so any imported profile is present when NM
# does its startup autoconnect scan.
Requires=glassos-system-dirs.service boot.mount
After=glassos-system-dirs.service boot.mount
RequiresMountsFor=/etc/NetworkManager/system-connections
Before=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/glassos-wifi-provision

[Install]
# multi-user.target places this in the same activation phase as NetworkManager.
# The Before=NetworkManager.service ordering above is only effective when both
# units share an activation transaction — network.target (the previous value) is
# reached *after* NM has already started, making Before= unreachable.
WantedBy=multi-user.target
