# buildroot-external/rootfs-overlay/usr/lib/systemd/system/glassos-bind.target
[Unit]
Description=GlassOS bind mounts
Documentation=man:systemd.target(5)
Before=network.target

[Install]
WantedBy=multi-user.target

