# buildroot-external/rootfs-overlay/usr/lib/systemd/system/NetworkManager.service.d/glassos.conf
[Unit]
# Ensure system directories exist and persistent paths are mounted before
# NetworkManager starts, so it can read persisted leases, state, and profiles.
Requires=glassos-system-dirs.service
After=glassos-system-dirs.service
RequiresMountsFor=/var/lib/NetworkManager /etc/NetworkManager/system-connections
